# linux-root/etc/udev/rules.d/99-drm.rules
ACTION=="add|change", SUBSYSTEM=="drm", RUN+="/opt/drm.sh"
