# linux-root/etc/udev/rules.d/99-keyboard.rules
# Check with lsusb or udevadm
ATTR{idVendor}=="04d9", ATTR{idProduct}=="4545", RUN+="/opt/keyboard.sh"
