# linux-root/etc/udev/rules.d/99-monitor.rules
ACTION=="add|change", SUBSYSTEM=="drm", RUN+="/opt/monitor.sh"
