# linux-root/etc/udev/rules.d/99-mouse.rules
# Check with lsusb or udevadm
ATTR{idVendor}=="046d", ATTR{idProduct}=="c069", RUN+="/opt/mouse.sh"
